# examples/systemd/binary/pours/deployment/signoz-telemetrystore-clickhouse-0-0.service
[Install]
WantedBy=multi-user.target

[Service]
CapabilityBoundingSet=CAP_NET_ADMIN CAP_IPC_LOCK CAP_SYS_NICE CAP_NET_BIND_SERVICE
ExecStart=/usr/bin/clickhouse-server --config=/etc/clickhouse-server/config.yaml --pid-file=/run/clickhouse/clickhouse.pid
Group=clickhouse
LimitCORE=infinity
LimitNOFILE=500000
Restart=always
RestartSec=30
RuntimeDirectory=clickhouse
Type=simple
User=clickhouse

[Unit]
After=time-sync.target network-online.target signoz-telemetrykeeper-clickhousekeeper-0.service
Description=ClickHouse Server (v25.5.6)
Requires=network-online.target
Wants=time-sync.target signoz-telemetrykeeper-clickhousekeeper-0.service
